# templates/make-snapshot.target
[Unit]
Description=Make Snapshot
Requires=multi-user.target
