# templates/rollback-snapshot.target
[Unit]
Description=Rollback Snapshot
Requires=multi-user.target
